# systemd/dbus-_0xdc.powersave.service.in
[Unit]
Description=Apply CPU frequency scaling and powersaving

[Service]
BusName=_0xdc.powersave
ExecStart=@@PREFIX@@/libexec/powersave/powersave.py

[Install]
WantedBy=multi-user.target
